# miner_lux.conf
[program:miner_lux]
command=/root/lux_mining.sh
process_name=%(program_name)s
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/lyngu/doraemon.log
