# miner_log.conf
[program:miner_log]
command=tailf /var/log/syslog > /dev/tty1
process_name=%(program_name)s
directory=/var/log/supervisor
user=root
autostart=true
autorestart=true
redirect_stderr=true
